# orderbook-bot.service
[Unit]
Description=Hyperliquid Trading Bot
After=network.target

[Service]
Type=simple
User=alexgd
WorkingDirectory=/home/alexgd/esilv/rust/orderbook-td
ExecStart=/home/alexgd/esilv/rust/orderbook-td/target/release/rust-3 trade
Restart=always
RestartSec=10
Environment=RUST_LOG=info
# Enable real trading
Environment=LIVE_TRADING=true

[Install]
WantedBy=multi-user.target
